# services/easychat-frontend.socket
[Unit]
Description=EasyChat Frontend Socket

[Socket]
Listen=/tmp/easychat-frontend.sock
Accept=no

[Install]
WantedBy=sockets.target